# water-monitor.service
[Unit]
Description=Water System Monitor
After=network.target auditd.service

[Service]
User=root
WorkingDirectory=/usr/local/plumbing-dashboard
ExecStart=/usr/local/plumbing-dashboard/water-monitor

[Install]
WantedBy=multi-user.target
Alias=water.service